# rootfs/etc/services.d/dashd-event-sock/run
#!/bin/sh

mkfifo /var/run/dogecoind/events-in
exec 8<>/var/run/dogecoind/events-in
exec ncat -l -k --send-only 0.0.0.0 22666 < /var/run/dogecoind/events-in
